# debian/prerm
#!/bin/sh

# Stop the service
systemctl stop gpt-for-uds.service

# Disable the service from starting on boot
systemctl disable gpt-for-uds.service

exit 0
